# docker/supervisor/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:translator]
command=python /app/translator/start_service.py
directory=/app/translator
user=meeshy
autostart=true
autorestart=true
startretries=3
startsecs=15
stdout_logfile=/app/logs/translator.log
stderr_logfile=/app/logs/translator.error.log
environment=DATABASE_URL="%(ENV_DATABASE_URL)s",REDIS_URL="%(ENV_REDIS_URL)s",GRPC_PORT="50051",HTTP_PORT="8000",ZMQ_PORT="5555",PYTHONUNBUFFERED="1"
priority=300

[program:gateway]
command=tsx /app/gateway/src/server.ts
directory=/app/gateway
user=meeshy
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/app/logs/gateway.log
stderr_logfile=/app/logs/gateway.error.log
environment=DATABASE_URL="%(ENV_DATABASE_URL)s",REDIS_URL="%(ENV_REDIS_URL)s",PORT="3000",NODE_ENV="production",ZMQ_TRANSLATOR_HOST="localhost",ZMQ_TRANSLATOR_PUSH_PORT="5555",ZMQ_TRANSLATOR_SUB_PORT="5558",ZMQ_TIMEOUT="30000"
priority=400

[program:frontend]
command=pnpm run dev
directory=/app/frontend
user=meeshy
autostart=true
autorestart=true
startretries=3
startsecs=15
stdout_logfile=/app/logs/frontend.log
stderr_logfile=/app/logs/frontend.error.log
environment=PORT="3100",NODE_ENV="development",NEXT_PUBLIC_API_URL="https://gate.meeshy.me",NEXT_PUBLIC_WS_URL="wss://gate.meeshy.me/"
priority=500

[group:meeshy]
programs=postgres,redis,translator,gateway,frontend,nginx
priority=999
